# n2k2ip.service
# systemd unit for the n2k2ip gateway (SocketCAN -> TCP/YDRAW).
#
# Install:
#   sudo apt install pipx                 # once, if you don't have it
#   sudo pipx install --global n2k2ip     # -> /usr/local/bin/n2k2ip
#   sudo cp n2k2ip.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now n2k2ip.service
#   journalctl -u n2k2ip.service -f
#
# Update later with: sudo pipx upgrade --global n2k2ip
#
# No-pip fallback: n2k2ip is a single stdlib-only file, so you can instead just
# copy it and point ExecStart at it:
#   sudo mkdir -p /opt/n2k2ip && sudo cp n2k2ip.py /opt/n2k2ip/
#   ExecStart=/usr/bin/python3 /opt/n2k2ip/n2k2ip.py --channel can0 --port 1457

[Unit]
Description=n2k2ip — NMEA2000 (SocketCAN) to TCP/YDRAW gateway
After=network.target

[Service]
Type=simple
User=root

# Bring can0 up at the NMEA2000 bitrate — but only if it isn't already up, so this
# is safe to run alongside another CAN service (e.g. a fastnet2n2k bridge) that
# shares can0: whoever starts first configures it, the others see it up and leave
# it alone (no interface flap, no start-up race). restart-ms 100 lets the
# controller auto-recover from a bus-off. The trailing check fails the unit if
# can0 can't be brought up at all.
ExecStartPre=/bin/sh -c 'ip link show can0 | grep -qw UP || ip link set can0 up type can bitrate 250000 restart-ms 100 2>/dev/null; ip link show can0 | grep -qw UP'

# /usr/local/bin/n2k2ip is where `pipx install --global` lands the command.
# Run `which n2k2ip` if yours differs.
ExecStart=/usr/local/bin/n2k2ip --channel can0 --port 1457
Restart=always
RestartSec=5

# === RESOURCE LIMITS ===
OOMScoreAdjust=-700
OOMPolicy=continue
MemoryMax=64M
MemoryHigh=48M
TimeoutStopSec=15

# === LOGGING ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n2k2ip
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
